# .config/systemd/user/conky2.service
[Unit]
Description=Conky, man

[Service]
Type=simple
ExecStart=/usr/bin/conky -b

[Install]
WantedBy=default.target
