# .config/systemd/user/udevedu.service
[Unit]
Description=Manage udev events in userspace

[Service]
Type=simple
ExecStart=/usr/bin/udevedu

[Install]
WantedBy=default.target
